# etc/profile.d/devtoolset-10.sh
#!/bin/sh

. /opt/rh/devtoolset-10/enable
