# etc/profile.d/irods.sh
#!/bin/sh

export spLogLevel=11
export IRODS_DATABASE_USER_PASSWORD_SALT=irodsPKey
