# etc/profile.d/llvm-toolset-7.sh
#!/bin/sh

. /opt/rh/llvm-toolset-7/enable
